# cron/crontab
*/15 * * * * /src/updater.sh 
*/15 * * * * cd /data/archived; python /src/organiser.py -e csv 
*/15 * * * * cd /data/archived; python /src/organiser.py -e TAB 
